# conf/supervisor.conf
[supervisord]

[supervisorctl]

[program:coinwatch]
command=python /app/src/coinwatch.py